# deploy/systemd/hanosuko-api.service
[Unit]
Description=Hanosuko.art API (views + tiers proxy)
After=network.target

[Service]
Type=simple
User=hanosuko
WorkingDirectory=/var/www/hanosuko/server
EnvironmentFile=/var/www/hanosuko/server/.env
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/hanosuko
PrivateTmp=true

[Install]
WantedBy=multi-user.target
