# deploy/systemd/link-profile-api.service
[Unit]
Description=Link Profile API (views + tiers proxy)
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/var/www/link-profile/server
EnvironmentFile=/var/www/link-profile/server/.env
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/link-profile
PrivateTmp=true

[Install]
WantedBy=multi-user.target
